# configs/cron.d/scan_worker-clearlog
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * */3 * * nico truncate -s 0 /tmp/pranger_scanworker.log >/dev/null 2>&1

